# admin/eix-sync.cron
#!/bin/env /bin/bash

/usr/bin/eix-sync
