# admin/qcheck.cron
#!/bin/env /bin/bash

MAIL=paperjam@localhost

qclog="/var/log/qcheck.${HOSTNAME}.${USER}.$(date +%y%m%d).log"

qcheck >> "${qclog}"
