# admin/qcheck.cron
#!/bin/env /bin/bash

MAILTO=paperjam@localhost

qcheck

#> /var/log/qcheck-"${HOSTNAME}"-"${USER}"-`date +%y%m%d%H%M%S`.txt
