# old/backup-paperjam.cron
#!/bin/env /bin/bash

/bin/tar -cvzf /mnt/ELEMENTS/linux/gentoo/genbox.paperjam.files.$(date +%y%m%d%H%M%S).tar.gz $(cat /home/paperjam/.backup.list.txt)
